# backup-utils/restic-backup/services/restic-backup.timer
[Unit]
Description=timer for system backup

[Timer]
OnUnitActiveSec=1d
OnBootSec=2m

[Install]
WantedBy=timers.target
